# debian/renasar-taskgraph.upstart
description "run the renasar taskgraph service"

env config=/etc/default/renasar-taskgraph

start on (started rabbitmq-server and started mongodb)
stop on runlevel [!2345]
respawn

pre-start script
[ -e "$config" ] || { echo "${config} missing, not starting"; stop; exit 0; }
end script

script
    echo Starting renasar taskgraph service
    chdir /var/renasar/renasar-taskgraph
    exec /usr/bin/nodejs /var/renasar/renasar-taskgraph/index.js
end script
